# Pwn/Ed-Edd-Eddie/crontab
* * * * * root tmpreaper 1m /tmp >> /var/log/cron.log 2>&1

